# crontab/mysql-backup-crontab
0 2 * * * bash /path-to/mysql-backup.sh
